# install/systemd/authbot.service.d/authbot.conf
[Service]
Environment="authbot_path=/home/user/authbot/"
